# config/systemd/techpulse-pulse.service
[Unit]
Description=TechPulse Full Intelligence Sequence
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/vishnu/worklab/techpulse
EnvironmentFile=/home/vishnu/worklab/techpulse/.env
ExecStart=/home/vishnu/.local/bin/uv run pulse run all
