# config/systemd/techpulse-collector.timer
[Unit]
Description=Run TechPulse AI Fast Collector Sweep every 3 hours

[Timer]
OnCalendar=*-*-* 00/3:00:00
Persistent=true

[Install]
WantedBy=timers.target
