# config/systemd/techpulse-pulse.timer
[Unit]
Description=Run TechPulse AI Full Intelligence Sequence every 6 hours

[Timer]
OnCalendar=*-*-* 00/6:00:00
Persistent=true

[Install]
WantedBy=timers.target
