# config/systemd/techpulse-collector.service
[Unit]
Description=TechPulse Fast Collector Sweep
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/vishnu/worklab/techpulse
EnvironmentFile=/home/vishnu/worklab/techpulse/.env
ExecStart=/home/vishnu/.local/bin/uv run pulse run collect
