# config/systemd/techpulse-keepalive.service
[Unit]
Description=Keep Supabase Database Alive
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/vishnu/worklab/techpulse
ExecStart=/home/vishnu/.local/bin/uv run python -c "print('Keeping workflows active')"
